# files/etc/hotplug.d/tty/10-vscom
#!/bin/sh
ln -s /dev/ttyO0 /dev/vscon
ln -s /dev/ttyO1 /dev/ttyVS0
ln -s /dev/ttyO2 /dev/ttyVS1

